# augy-music-bot.service
[Unit]
Description=HelloTunes
After=network.target lavalink.service
Wants=lavalink.service

[Service]
Type=simple
User=augy
Group=augy
WorkingDirectory=/home/augy/hellotunes
ExecStart=/home/augy/hellotunes/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hellotunes

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
